# macos/LaunchAgents/com.simon.sticks3-codex-companion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.simon.sticks3-codex-companion</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/simon/Documents/workspace/repos/sticks3-codex-companion/scripts/sticks3-macos-bridge</string>
    <string>run</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/simon/Documents/workspace/repos/sticks3-codex-companion</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>STICKS3_REPO</key>
    <string>/Users/simon/Documents/workspace/repos/sticks3-codex-companion</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/sticks3-codex-companion.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/sticks3-codex-companion.launchd.err.log</string>
</dict>
</plist>
